# packages/cron/etc/cron.d/1create-topics
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root
HOME=/
*/2 * * * * root /opt/bookit/bin/create-topics-after-startup.sh
